# files/system/usr/lib/systemd/system/hp-set-lighting.service
[Unit]
Description=Load HP Module to set HP Lighting on Boot

[Service]
ExecStart=/usr/bin/hp-set-lighting

[Install]
WantedBy=multi-user.target